# rpm/kubelet.service
[Unit]
Description=Kubernetes Kubelet Server
Documentation=https://github.com/kubernetes/kubernetes

[Service]
EnvironmentFile=/etc/sysconfig/kubelet
ExecStart=/var/lib/kubelet/kubelet-wrapper
Restart=always

[Install]
WantedBy=multi-user.target

